# data/gnome-user-share-webdav.service.in
[Unit]
Description=gnome-user-share WebDAV server

[Service]
Type=dbus
BusName=org.gnome.user-share.webdav
ExecStart=@libexecdir@/gnome-user-share-webdav
Restart=on-failure
